# backend/paper_trading/supervisor.conf
[program:paper_trading]
command=/root/.venv/bin/python3 /app/backend/paper_trading/service.py
directory=/app/backend/paper_trading
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/paper_trading.err.log
stdout_logfile=/var/log/supervisor/paper_trading.out.log
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
